# .config/systemd/user/autocutsel.service
[Unit]
Description=autocutsel

[Service]
ExecStart=/usr/bin/autocutsel -selection CLIPBOARD

[Install]
WantedBy=default.target
